# systemd/user/pipewire.socket.d/override.conf
[Socket]
ExecStartPost=/bin/mount /srv/win10/pipewire-0
ExecStopPre=/bin/umount /srv/win10/pipewire-0
